# units/email-status@.service.d/recipient-email.conf.j2
[Service]
Environment='RECIPIENT_EMAIL={{ recipient_email }}'
